# net-dns/mini-dot/files/mini-dot.openrc
#!/sbin/openrc-run
# Copyright 1999-2020 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

command="/usr/bin/forwarder.py"
command_args="${START_OPTS}"
pidfile="/run/${RC_SVCNAME}.pid"
command_background=true
#command_user="fpnd:fpnd"

depend() {
	provide dns
	need localmount net
}
