# net-misc/zerotier/files/zerotier.service
[Unit]
Description=ZeroTier One
After=network.target
Wants=network-online.target
PartOf=fpnd.service

[Service]
ExecStart=/usr/sbin/zerotier-one/zerotier-one
Restart=always
KillMode=process
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
